# scripts/guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zhuchenyuan.congxicai-v6.guardian</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/zhuchenyuan/工作流/cong-xi-fa-cai-v6/scripts/guardian.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/zhuchenyuan/工作流/cong-xi-fa-cai-v6</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/zhuchenyuan/cong-xi-fa-cai-logs/guardian.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/zhuchenyuan/cong-xi-fa-cai-logs/guardian.stderr.log</string>
</dict>
</plist>
